# webhook-server.service
[Unit]
Description=GitHub Webhook Server for eBay Auto Parts Lister
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/ebay-auto-parts-lister
ExecStart=/opt/ebay-auto-parts-lister/venv/bin/python3 /opt/ebay-auto-parts-lister/webhook-server.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment
Environment=PYTHONPATH=/opt/ebay-auto-parts-lister

[Install]
WantedBy=multi-user.target
